# crontab
* * * * * root /app/load_playlist.sh >> /var/log/cron_output.log 2>&1
