# crontab
* * * * * root /usr/local/bin/python3 /app/load_playlist.py >> /var/log/cron_output.log 2>&1
